# root/etc/services.d/deluge-web/run
#!/bin/sh

exec s6-setuidgid deluge \
       /usr/bin/deluge-web \
         --config /var/lib/deluge \
         --loglevel info
